# deploy/systemd/aipg-horde@.service
# Flask (legacy queue) worker — one instance per port.
# Enable the fleet:  systemctl enable --now aipg-horde@{7001..7008}
[Unit]
Description=AIPG Horde Flask API (port %i)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=aipg
Group=aipg
WorkingDirectory=/home/aipg/system-core
EnvironmentFile=/etc/aipg/grid.env
ExecStart=/home/aipg/system-core/.venv/bin/python server.py -vv --horde stable -p %i
Restart=always
RestartSec=3
# Logs go to journald: journalctl -u aipg-horde@7001 -f
StandardOutput=journal
StandardError=journal
# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
